# cpu-monitoring.service
[Unit]
Description=CPU Monitoring Service

[Service]
ExecStart=/bin/bash /usr/bin/cpu-monitoring

[Install]
WantedBy=multi-user.target